# unbound-monthly.cron
#!/bin/sh

/usr/sbin/unbound-anchor -a /etc/unbound/root.anchor -c /etc/unbound/icannbundle.pem
