# unbound-monthly.cron
#!/bin/sh

/usr/sbin/unbound-anchor -a /var/lib/unbound/root.key -c /etc/unbound/icannbundle.pem
